# WaarpVitamInstall/waarp/conf/r66/systemd/system/waarp-vitam-ingest.service
[Unit]
Description=Waarp-Vitam-Ingest-Monitor
Documentation=http://www.waarp.fr
Before=systemd-user-sessions.service
After=network.target
After=syslog.target

[Service]
Type=simple
User=vitam
Group=vitam
ExecStart=%ROOT%/bin/r66/IngestMonitor.sh start

[Install]
WantedBy=multi-user.target
